# .config/systemd/user/hypr-sun.service
[Unit]
Description=Update Hyprpaper Wallpaper based on Sun Position
# Wir entfernen den harten Ordering-Cycle
After=graphical-session.target
# Erlaubt dem System 5 Versuche innerhalb von 60s
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=oneshot
# DER FIX: Wir zwingen systemd zu einer Pause VOR dem Skriptstart,
# ohne das Skript selbst zu verändern.
ExecStartPre=/usr/bin/sleep 3
ExecStart=%h/.local/bin/hypr-sun.sh

# Error-Handling wie gehabt, aber mit Puffer
Restart=on-failure
RestartSec=1s

[Install]
# Wichtig: Bindung an das User-Target für stabilere Umgebung
WantedBy=default.target
